# deploy/systemd/alloc-context-weekly-brief.timer
[Unit]
Description=AllocContext weekly brief timer
Requires=alloc-context-weekly-brief.service

[Timer]
OnCalendar=Mon *-*-* 07:00:00 America/Chicago

[Install]
WantedBy=timers.target
